# scripts/host/jarvis-kiosk-restart.service
[Unit]
Description=Reinicio nocturno del kiosko J.A.R.V.I.S. (libera memoria de Chromium)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart jarvis-kiosk.service
